# net-wireless/intel-wimax-network-service/files/wimax
#!/sbin/runscript

description="WiMAX Network Service"

depend() {
	need localmount
	after bootmisc
}

start() {
	ebegin "Starting WiMAX Network Service"
	[ -e /var/run/wimaxd.pid ] && rm -f /var/run/wimaxd.pid
	start-stop-daemon --start --quiet --pidfile /var/run/wimax.pid --exec /usr/bin/wimaxd
	eend $?
}

stop() {
	ebegin "Stoping WiMAX Network Service"
	start-stop-daemon --stop --quiet --pidfile /var/run/wimax.pid
	eend $?
}
