# x11-misc/mkxf86config/files/mkxf86config.initd
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="Xorg-x11 configuration builder"

depend() {
	before local xdm
}

start() {
	ebegin "Creating X Configuration"
    /usr/sbin/mkxf86config.sh
    eend $?
}
