# x11-misc/mkxf86config/files/mkxf86config.initd
#!/sbin/runscript
 
depend() {
    before local xdm
}
 
start() {
    ebegin "Creating X Configuration"
    /usr/sbin/mkxf86config.sh
	eend $?
}
