# www-apps/seafile-pro-server/files/seafile-server
#!/sbin/openrc-run

#command=/usr/bin/foo
#command_args="${foo_args} --bar"
#pidfile=/var/run/seafile-server.pid

name="Seafile Daemon"
description="Seafile server"

#change these
seafile_dir="<CHANGE ME>"
seafile_user="seafile"

start() {
    cd ${seafile_dir}
    su -c "${seafile_dir}/seafile-server-latest/seafile.sh start >/dev/null 2>/dev/null" ${seafile_user}
    su -c "${seafile_dir}/seafile-server-latest/seahub.sh start >/dev/null 2>/dev/null" ${seafile_user}
}

stop() {
    cd ${seafile_dir}
    su -c "${seafile_dir}/seafile-server-latest/seahub.sh stop  >/dev/null 2>/dev/null" ${seafile_user}
    su -c "${seafile_dir}/seafile-server-latest/seafile.sh stop >/dev/null 2>/dev/null" ${seafile_user}
}
