# sys-power/acpid/files/systemd/acpid.socket
[Unit]
Description=ACPI Event Daemon Listen Socket

[Socket]
ListenStream=/run/acpid.socket

[Install]
WantedBy=sockets.target
