# net-vpn/tuntox/files/tuntox.conf
[program:tuntox]
command = /usr/sbin/tuntox
autostart = true
startsecs = 20
user = tuntox
redirect_stderr = true
stdout_logfile = /var/log/tuntox/tuntox.log
