# sys-power/acpid/files/acpid.service
[Unit]
Description=ACPI Event Daemon
After=syslog.target

[Service]
ExecStart=/usr/sbin/acpid -f

[Install]
WantedBy=multi-user.target
